# .config/systemd/user/alarmclock.service
[Unit]
Description=Alarm Clock

[Service]
ExecStart=/usr/bin/mpv --no-terminal --loop-file=5 /home/diego/Downloads/sedna_ringtone.mp3
